# resources/autohotspot/autohotspot.timer
[Unit]
Description=autohotspot timer

[Timer]
OnActiveSec=5min
OnUnitInactiveSec=5min

[Install]
WantedBy=basic.target
